# console-hub.service
[Unit]
Description=Golden Tide Console Hub
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/YOUR_USER/golden-tide-console/console-hub
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
Environment=PORT=3000 NODE_ENV=production

[Install]
WantedBy=multi-user.target
